# feeds/clsgui/modules/clsgui-mod-wireless/root/etc/uci-defaults/50_clsgui-mod-wireless
#!/bin/sh

uci -q batch <<-EOF >/dev/null
	set luci.wifi_dualband=internal
	commit luci
EOF

exit 0
