# target/linux/clsemi/base-files/etc/uci-defaults/10_firewall
#!/bin/sh

uci -q get firewall.@defaults[0].flow_offloading > /dev/null && exit 0

uci set firewall.@defaults[0].flow_offloading=1
uci set firewall.@defaults[0].flow_offloading_hw=1
uci commit firewall

exit 0
